# init/launchd/com.elodinlaarz.agent-routines.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.elodinlaarz.agent-routines</string>
  <key>ProgramArguments</key>
  <array>
    <!-- Replace with the absolute path to `routines`. -->
    <string>/usr/local/bin/routines</string>
    <string>daemon</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/tmp/agent-routines.log</string>
  <key>StandardErrorPath</key><string>/tmp/agent-routines.log</string>
</dict>
</plist>
